# agente73/deploy/agente73.conf
; Supervisor — Agente 73 v4.0 "Con Alma"
; Instalar en /etc/supervisor/conf.d/agente73.conf (ver RUNBOOK.md)
[program:agente73]
command=/opt/agente73/venv/bin/gunicorn --bind 127.0.0.1:8091 --workers 2 --threads 4 --timeout 120 --access-logfile /var/log/agente73/access.log --error-logfile /var/log/agente73/error.log agente73.app:app
directory=/opt/agente73
user=agente73
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
environment=HOME="/home/agente73",PATH="/opt/agente73/venv/bin:/usr/local/bin:/usr/bin:/bin"
stdout_logfile=/var/log/agente73/supervisor.log
stderr_logfile=/var/log/agente73/supervisor.err
